# user/lede-openclash-x86-amd64/files/etc/uci-defaults/zzzzz-extra-settings
#!/bin/sh

# openclash
uci set openclash.config.enable_redirect_dns='0'

uci commit openclash

uci commit

exit 0
