# user/immortalwrt-openclash-x86-amd64/files/etc/uci-defaults/zzzzz-extra-settings
#!/bin/sh

uci set luci.main.mediaurlbase=/luci-static/argon
uci commit luci

exit 0
